# battery-monitor.service
[Unit]
Description=Battery Monitor
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/battery-monitor
Restart=on-failure
RestartSec=30s
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

[Install]
WantedBy=default.target
